# deploy/website-deploy.service
[Unit]
Description=Website deploy webhook
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=codex
Group=codex
EnvironmentFile=/etc/website-deploy.env
WorkingDirectory=/opt/website/app
ExecStart=/usr/bin/node /opt/website/app/deploy/webhook-server.mjs
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
